# dexter-installer_1.0_all/DEBIAN/prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit/13.6ubuntu1
if [ -z "${DPKG_ROOT:-}" ] && [ "$1" = remove ] && [ -x "/etc/init.d/live-installer" ] ; then
	invoke-rc.d live-installer stop || exit 1
fi
# End automatically added section
# Automatically added by dh_systemd_start/13.6ubuntu1
if [ -z "${DPKG_ROOT:-}" ] && [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	deb-systemd-invoke stop 'install-live-drivers.service' >/dev/null || true
fi
# End automatically added section
